# setup/iheart-radio.service
[Unit]
Description=Oasis Radio Stream Deck Controller
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nashville-radio/streamdeck/controller.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
